# deploy/fantasyyc-metadata.service
[Unit]
Description=FantasyYC NFT Metadata Server (Express.js on port 3001)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=fantasyyc
Group=fantasyyc
WorkingDirectory=/opt/fantasyyc/backend
ExecStart=/usr/bin/node server.js
EnvironmentFile=/opt/fantasyyc/.env
Environment=NODE_ENV=production
Environment=PORT=3001

# Restart policy
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

KillSignal=SIGTERM
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/fantasyyc/backend
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fantasyyc-metadata

[Install]
WantedBy=multi-user.target
